# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Use runtipi-cli prepare to set up directory structure and configuration
        # This matches the official runtipi installation process
        cd /opt/runtipi && ./runtipi-cli prepare

        # Reload systemd daemon to pick up the new service
        systemctl daemon-reload

        # Enable the service to start on boot
        systemctl enable runtipi.service

        # Start the service
        systemctl start runtipi.service
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
